# monitor/currentcost/01UsbAdded.rules
KERNEL=="ttyUSB*",  ACTION=="add", RUN+="/usr/bin/currentcost-monitor %r/%k"
KERNEL=="ttyUSB*",  ACTION=="remove", RUN+="/usr/bin/currentcost-stop"
